# tests_and_snippets/sh_and_py/get_version/test_cases/systemd_unit_versionline.service
[Unit]
Description=WsprryPi Service
Version=1.2.1-8c2314a-dirty [version_exe]

[Service]
ExecStart=/usr/bin/wsprrypi

[Install]
WantedBy=multi-user.target
